# stow/systemd/.config/systemd/user/cachyos-wiki-pull.service
[Unit]
Description=Pull latest CachyOS wiki mirror (~/AI/Docs/cachyos-wiki)

[Service]
Type=oneshot
Nice=10
ExecStart=/usr/bin/bash %h/AI/Docs/cachyos-wiki-pull.sh
