# stow/systemd/.config/systemd/user/cachyos-wiki-pull.timer
[Unit]
Description=Monthly pull of the local CachyOS wiki mirror

[Timer]
OnCalendar=monthly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
