# stow/systemd/.config/systemd/user/niri-monitor-toggle.service
[Unit]
Description=Niri monitor auto switch

[Service]
Type=oneshot
ExecStart=%h/.local/bin/niri-monitor-toggle.sh
