# stow/systemd/.config/systemd/user/openclaw-node.service
[Unit]
Description=OpenClaw Node Host (v2026.7.1)
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
ExecStart=/usr/bin/node /home/zefuros/.local/share/npm/lib/node_modules/openclaw/dist/index.js node run --host openclaw.zecuritylabs.com --port 443 --tls
Restart=always
RestartSec=5
RestartPreventExitStatus=78
TimeoutStopSec=30
TimeoutStartSec=30
SuccessExitStatus=0 143
OOMPolicy=continue
KillMode=control-group
EnvironmentFile=-/home/zefuros/AI/Harnesses/OpenClaw/state/node.systemd.env
Environment=HOME=/home/zefuros
Environment=TMPDIR=/tmp
Environment=OPENCLAW_STATE_DIR=/home/zefuros/AI/Harnesses/OpenClaw/state
Environment=OPENCLAW_CONFIG_PATH=/home/zefuros/AI/Harnesses/OpenClaw/state/openclaw.json
Environment=PATH=/usr/bin:/usr/local/bin:/bin:/home/zefuros/.local/share/npm/bin:/home/zefuros/.local/bin:/home/zefuros/.npm-global/bin:/home/zefuros/bin:/home/zefuros/.nix-profile/bin
Environment=OPENCLAW_LAUNCHD_LABEL=ai.openclaw.node
Environment=OPENCLAW_SYSTEMD_UNIT=openclaw-node
Environment="OPENCLAW_WINDOWS_TASK_NAME=OpenClaw Node"
Environment=OPENCLAW_WINDOWS_TASK_HIDDEN_LAUNCHER=1
Environment=OPENCLAW_TASK_SCRIPT_NAME=node.cmd
Environment=OPENCLAW_LOG_PREFIX=node
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=node
Environment=OPENCLAW_SERVICE_VERSION=2026.7.1

[Install]
WantedBy=default.target
