# stow/systemd/.config/systemd/user/prox-weekly-sweep.timer
[Unit]
Description=Run the Prox weekly homelab sweep every Monday morning

[Timer]
OnCalendar=Mon *-*-* 08:07
Persistent=true

[Install]
WantedBy=timers.target
