# stow/systemd/.config/systemd/user/rustdesk.service
[Unit]
Description=RustDesk Remote Desktop

[Service]
ExecStart=/usr/bin/rustdesk
Restart=on-failure

[Install]
WantedBy=default.target
